# templates/db-backup.cron.j2
50 1 * * * root /usr/local/bin/db-backup.sh {{mode}} {{database}}
