# templates/s3-mysql-backup.cron.j2
50 1 * * * root s3-mysql-backup {{database}} /etc/s3-mysql-backup/{{database}}-backup.yml
