# templates/s3-mongo-backup.cron.j2
30 1 * * * root mongodb_s3_backup -n /etc/s3-mongo-backup/{{database}}-backup.json
