# files/haproxy-supervisor.conf
[program:haproxy]
command=/usr/sbin/haproxy -f /etc/haproxy/haproxy.cfg -db
autorestart=true
redirect_stderr=true
