# files/logclean.cron
# Delete old log files
17 3      * * *    root   find /var/log/app -type f -mtime +7 -delete
28 3      * * *    root   find /var/log/sys -type f -mtime +7 -delete
