# files/mongodb-supervisor.conf
[program:mongodb]
command=/usr/bin/mongod --config /etc/mongodb.conf
autorestart=true
redirect_stderr=true
