# files/mongodb-supervisor.conf
[program:mongodb]
command=/usr/sbin/mongod
autorestart=true
redirect_stderr=true
