# files/mysql-supervisor.conf
[program:mysql]
command=/usr/sbin/mysqld
autorestart=true
redirect_stderr=true
