# files/nginx-supervisor.conf
[program:nginx]
command=/usr/sbin/nginx
autorestart=true
redirect_stderr=true
startsecs=5
