# files/php5-fpm-supervisor.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm
autorestart=true
redirect_stderr=true
