# files/rabbitmq-supervisor.conf
[program:rabbitmq]
command=/usr/lib/rabbitmq/bin/rabbitmq-server
autorestart=true
redirect_stderr=true
environment=HOME=/tmp
